# roles/kurtzhong.glide/templates/etc/profile.d/glide.sh
export PATH=$PATH:/opt/glide/{{ package_os }}-{{ package_arch }}
